# services/alsa-grpc-server/finish
#!/bin/bash
echo 'alsa-grpc-server exited'