# services/cpiped/run
#!/usr/bin/with-contenv /bin/bash
source /root/dynamic.env
echo "Starting cpiped ${HW_ADDRESS} > ${PIPE}"
exec /cpiped/cpiped -d ${HW_ADDRESS} ${PIPE}

